# deploy/wvs-streamlit.service
[Unit]
Description=WVS Bot Streamlit UI
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/python/wvs_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/python/wvs_bot/.venv/bin/streamlit run ui/streamlit_app.py \
  --server.address=127.0.0.1 \
  --server.port=8502 \
  --server.headless=true
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
